# embedded-linux/K26.1/nguyen-sy-hung/session-06/Exercise_2/logger.service
[Unit]
Description=Embedded Linux Logger Service
After=multi-user.target


[Service]
Type=simple

ExecStart=/usr/local/bin/logger

StandardOutput=journal
StandardError=journal

Restart=no


[Install]
WantedBy=multi-user.target