# embedded-linux/K26.1/cao-duc-anh/session-06/Exercise_2/logger.service
[Unit]
Description=Logger crash-test service

[Service]
Type=simple
ExecStart=/usr/local/bin/logger
StandardOutput=journal
StandardError=journal
Restart=no
