# embedded-linux/K26.1/bui-quang-huy/session-06/Exercise_2/logger.service
[Unit]
Description=Logger Service for Severity Level Analysis
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/logger
StandardOutput=journal
StandardError=journal
Restart=no

[Install]
WantedBy=multi-user.target