# embedded-linux/K26.1/nguyen-quoc-viet/session-06/Exercise_2/logger.service
[Unit]
Description=Logger Service — severity-level demo
After=network.target

[Service]
ExecStart=/usr/local/bin/logger_demo
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
