# embedded-linux/K26.1/cao-duc-anh/session-06/Exercise_1/monitor.service
[Unit]
Description=Monitor background service
After=network.target
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/local/bin/monitor
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
